# common/task/.config/systemd/user/task-sync.service
[Unit]
Description=Synchronize Taskwarrior server

[Service]
type=oneshot
ExecStart=/usr/bin/task sync
