# common/task/Library/LaunchAgents/com.patrickanker.twgcalsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.patrickanker.twgcalsync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/patrickanker/.local/bin/twgcalsync.sh</string>
    </array>

    <key>Nice</key>
    <integer>1</integer>

    <key>StartInterval</key>
    <integer>5</integer>

    <key>StandardErrorPath</key>
    <string>/tmp/twgcalsync.err</string>

    <key>StandardOutPath</key>
    <string>/tmp/twgcalsync.out</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
