# common/kde/.config/autostart/ssh-add.sh.desktop
[Desktop Entry]
Exec=/home/psanker/.local/bin/ssh-add.sh
Icon=dialog-scripts
Name=ssh-add.sh
Type=Application
X-KDE-AutostartScript=true
